# conf/supervisor/supervisord.conf
[unix_http_server]
file=/run/supervisord.sock   ; (the path to the socket file)

[supervisord]
logfile=/var/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
loglevel=info                ; (log level;default info; others: debug,warn,trace)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock ; use a unix:// URL  for a unix socket

[program:api-service]
command=uwsgi -psgi /etc/uwsgi/uwsgi.ini --enable-threads
# command=python manage.py run
directory=/webapps
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/apps/api.log

[program:celery-service]
command=celery -A manage.celery_app worker -l INFO -c 2 -n celery@%%h --queues=celery -E
directory=/webapps
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/apps/celery.log
